# services/linux/pact.socket
[Unit]
Description=PACT gunicorn socket

[Socket]
ListenStream=/run/pact-gunicorn.socket

[Install]
WantedBy=sockets.target
